# roles/openstack_block-storage/templates/rc.local.j2
#!/bin/bash

ip link set br-ex up
ip link set {{ ens0 }} up
ip addr add {{ local_ip }}/24 dev br-ex
ip route add default via {{ gateway }}

losetup /dev/loop2 /dev/zd0
exit 0
